# configs/systemd/user/snry-daemon.service
[Unit]
Description=Snry Shell Daemon
Documentation=https://github.com/sonroyaalmerol/snry-shell
After=graphical-session.target hyprland-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/snry-daemon daemon
Restart=on-failure
RestartSec=5

# Security hardening
ProtectSystem=strict
ReadWritePaths=%h/.config %h/.cache %h/.local %h/.local/state
PrivateTmp=true
NoNewPrivileges=true

# Environment
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-1

[Install]
WantedBy=graphical-session.target
