# configs/systemd/user/snry-daemon.service
[Unit]
Description=Snry Shell Daemon
Documentation=https://github.com/sonroyaalmerol/snry-shell
# No After=graphical-session.target — the daemon launches Hyprland itself.
# It should start early in the user session.

[Service]
Type=simple
ExecStart=/usr/bin/snry-daemon daemon
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target